# dockerfiles/itest/itest/nerve.conf
start on runlevel [2,3,4,5]
description "nerve"
script
  exec setuidgid nobody nerve --config /etc/nerve/nerve.conf.json
end script
